# etc/systemd/system/memlockd.service.d/override.conf
[Unit]
ConditionPathExists=/etc/memlockd.d/vanished_root_handler.cfg
